# systemd/rhythm-jump.service
[Unit]
Description=Rhythm Jump Backend Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/rhythm-jump/backend
Environment=PATH=/usr/local/bin:/usr/bin:/bin
Environment=RHYTHM_HEADLESS_MODE=1
ExecStart=/bin/sh -lc 'set -eu; UV_BIN="$(command -v uv || true)"; if [ -z "$UV_BIN" ]; then echo "uv not found in PATH: $PATH" >&2; exit 127; fi; exec "$UV_BIN" run --project /opt/rhythm-jump/backend uvicorn rhythm_jump.main:app --host 0.0.0.0 --port 8000'
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
